# assets/etc/supervisord/license-server.conf
[program:license-server]
command = /usr/local/bin/license-server run
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
